# h1a_hp.service
[Unit]
Description=H1A Hotplug Utility service

[Service]
Type=simple
ExecStart=/usr/local/sbin/h1a_hp

[Install]
WantedBy=multi-user.target
